# config/boot-notify.service
[Unit]
Description=Boot notification email (IP locale via Resend)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=minitel
WorkingDirectory=/home/minitel/minitel-gpt
EnvironmentFile=/home/minitel/minitel-gpt/.env
ExecStart=/usr/bin/python3 /home/minitel/minitel-gpt/services/boot_notify.py
StandardOutput=append:/home/minitel/minitel-gpt/logs/notify.log
StandardError=append:/home/minitel/minitel-gpt/logs/notify.log

[Install]
WantedBy=multi-user.target
