# config/wifi-manager.service
[Unit]
Description=WiFi Manager — provisioning hotspot style IoT
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/minitel/minitel-gpt
EnvironmentFile=/home/minitel/minitel-gpt/.env
ExecStart=/usr/bin/python3 /home/minitel/minitel-gpt/services/wifi_manager.py
Restart=always
RestartSec=10
StandardOutput=append:/home/minitel/minitel-gpt/logs/wifi.log
StandardError=append:/home/minitel/minitel-gpt/logs/wifi.log

[Install]
WantedBy=multi-user.target
